# udev/060f_ic0072_slot22.rules

ACTION=="add|remove", SUBSYSTEM=="tty", DEVPATH=="*/tty/ttymxc3", RUN+="/opt/unipi/tools/uhelper rs232.22 tty 22 $name"