# udev/060f_ic0072_slot21.rules

ACTION=="add|remove", SUBSYSTEM=="tty", DEVPATH=="*/tty/ttymxc1", RUN+="/opt/unipi/tools/uhelper rs232.21 tty 21 $name"