# udev/ic0072_slot12.rules

ACTION=="add|remove", SUBSYSTEM=="tty", DEVPATH=="*/tty/ttymxc0", RUN+="/opt/unipi/tools/uhelper rs232.12 tty 12 $name"